# ods.rules
SUBSYSTEM=="block", ENV{ID_CDROM}=="?*", ENV{ID_SERIAL}=="ATAPI_DVD_A_DH16ABSH_KU0160F01110904D9DF3", ACTION=="change", RUN+="/opt/ods.sh"
